# profile.d/personal_scripts.sh
export PATH=$PATH:/home/georgi/.local/bin/scripts
